# target/linux/mvebu64/base-files/etc/uci-defaults/02-rps
#!/bin/sh
#
# Copyright (C) 2017 OpenWrt.org
#

[ ! -e /etc/config/network ] && exit 0

# Disable RPS by default
[[ -z $(uci -q get network.globals.default_ps) ]] || [[ $(uci get network.globals.default_ps) == '1' ]] && $(uci set network.globals.default_ps=0)

uci commit network

exit 0
