# package/network/services/hostapd/files/hostapd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2017 OpenWrt.org

START=90

USE_PROCD=1
PROG=/usr/sbin/hostapd
CONFIG=/etc/hostapd.conf


start_service () {
	[[ ! -f /etc/hostapd.conf ]] && CONFIG=/rom/etc/hostapd.conf

	$PROG -B $CONFIG
}
stop_service() {

	killall hostapd
	sleep 2
}
